# logban.cron
*/5	*	*	*	*	. /etc/logban.conf; /bin/sh ${logbandir}/cklog.sh
